# src/etc/cron.d/wlan-monit
# /etc/cron.d/network-monitor
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Network connectivity monitor
# Every 15 minutes, offset +7 min from quarter-hour to avoid colliding with
# piSolar Renogy RS485 USB reads (*/5 on the 5-minute marks).
7,22,37,52 * * * * root /usr/local/bin/wlan-monit >/dev/null 2>&1
