# src/usr/lib/systemd/system/wol.service
[Unit]
Description=Configure Wake-up on LAN

[Service]
Type=oneshot
ExecStart=/usr/sbin/ethtool -s _ETH_DEV_ wol g

[Install]
WantedBy=basic.target
